# ml-broker/ml-broker.service.in
[Unit]
Description=ML message broker
Wants=network.target
After=network.target
After=redis.service

[Service]
Type=simple
RestartSec=1s
Restart=always
ExecStart=python3 /opt/ml-tools/ml-broker/ml-broker.py
WorkingDirectory=/opt/ml-tools/ml-broker

[Install]
WantedBy=multi-user.target
